# deploy/study-plan-backup.service
[Unit]
Description=Back up Summer Study Plan data

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backup-study-plan
