# deploy/study-plan-backup.timer
[Unit]
Description=Daily Summer Study Plan backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true

[Install]
WantedBy=timers.target
